# kill/task.cron
* * * * * touch /home/zhu/shellscript/kill/a.file
